# claude-architect.service
[Unit]
Description=Claude Architect Service

[Service]
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=%h/.local/bin/claude-architect
Restart=on-failure

[Install]
WantedBy=default.target
